# systemd/fortheweebs.service
# Systemd service for non-K8s deployments (bare metal, VMs)
[Unit]
Description=lunonex API Server
After=network.target

[Service]
Type=simple
User=lunonex
WorkingDirectory=/opt/lunonex
ExecStart=/usr/bin/node server.js

# Auto-restart on crash
Restart=always
RestartSec=5

# Disable start limits (avoid permanent backoff)
StartLimitIntervalSec=0
StartLimitBurst=0

# Environment
Environment="NODE_ENV=production"
Environment="PORT=3001"
EnvironmentFile=/etc/lunonex/env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lunonex-api

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
MemoryLimit=2G

[Install]
WantedBy=multi-user.target
